# database/settingan_cron_backup.cron
0 * * * * bash /home/alvancho/Documents/IO5/FP/client/program_dumb_client